# scripts/power_unplugged.rules
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/etc/scripts/buptime/on_unplugged.sh"
